# network/gost/crontab
0 0 * * * /usr/local/bin/certbot-auto renew >/dev/null
0 0 * * * rm /root/root.log
